# scripts/fix-permissions.cron
0 2 * * * bash /home/flip/tars-ai-project/scripts/fix-permissions.sh
